# aichat/scripts/aichat.service
[Unit]
Description=AI Chat Application - Trợ lý AI thông minh
After=network.target graphical.target
Wants=graphical.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/aichat
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"
Environment="QT_QPA_PLATFORM=xcb"
ExecStart=/home/pi/aichat/scripts/start_aichat.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Giới hạn tài nguyên (tùy chọn)
# MemoryMax=512M
# CPUQuota=60%

[Install]
WantedBy=graphical.target

